# service/cube-network-monitor.service
[Unit]
Description=Cube Network Monitor

[Service]
ExecStart=/data/cube/scripts/network-monitor.sh
User=root
Restart=on-abort

[Install]
WantedBy=multi-user.target
WantedBy=network-online.target
